# src/pkp-ojs/assets/crontab
0,30 * * * * sh /usr/local/bin/run-scheduler.sh
